# config/includes.chroot/etc/rc.local
#!/bin/bash
# Disable the network devices
/sbin/ifconfig eth0 down
/sbin/ifconfig wlan0 down
# Spoof the mac addresses
/usr/bin/macchanger -r eth0
/usr/bin/macchanger -r wlan0
# Re-enable the Ethernet device only, the wifi stays down.
/sbin/ifconfig eth0 up
#ifconfig wlan0 up

/bin/bash -c '. /usr/bin/fireblock && fireblock_watch_xscreensaver &'
/bin/bash -c '. /usr/bin/hidblock && hidblock_watch_xscreensaver &'

exit 0